# docker/config/supervisord/supervisord_app.conf
[unix_http_server]
file=/run/supervisord.sock
username = www-data
password = www-data

[supervisord]
pidfile=/run/supervisord.pid
loglevel=warn
logfile_maxbytes=10MB  
logfile=/run/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=php-fpm
autostart=true
autorestart=true
numprocs=8
user=www-data